# deploy/ubuntu/bot-whatsapp.service
[Unit]
Description=Travaux WhatsApp Bot
After=network.target mysql.service

[Service]
Type=simple
User=travaux
Group=travaux
WorkingDirectory=/home/travaux/app/Travaux/Travaux/Bot
Environment=NODE_ENV=production
Environment=BOT_ENV_FILE=/home/travaux/.travaux/bot.env
Environment=BOT_AUTH_DIR=/home/travaux/.travaux/bot/auth_info_baileys
Environment=BOT_FORCE_NEW_QR=1
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full

[Install]
WantedBy=multi-user.target